# deploy/com.secondbrain.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd template for the SecondBrain capture+transcribe daemon.

  1. Edit the paths below (REPLACE __REPO__ and __PYTHON__).
  2. Copy to ~/Library/LaunchAgents/com.secondbrain.daemon.plist
  3. launchctl load ~/Library/LaunchAgents/com.secondbrain.daemon.plist
     (macOS will prompt for Microphone permission the first time.)

  A second, similar agent can run `sb serve` for the web UI, and a third can run
  `python -m secondbrain.menubar.app` for the menu bar control.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.secondbrain.daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON__</string>
        <string>-m</string>
        <string>secondbrain.daemon</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__REPO__</string>

    <!-- launchd starts agents with a minimal PATH that excludes Homebrew, so the
         daemon can't find ffmpeg (transcription decodes FLAC via ffmpeg) and every
         transcribe job fails with "FFmpeg is not installed or not in your PATH."
         Include the Homebrew bin dirs for Apple Silicon (/opt/homebrew) and
         Intel (/usr/local). -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__REPO__/data/daemon.out.log</string>
    <key>StandardErrorPath</key>
    <string>__REPO__/data/daemon.err.log</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
